# daemon/com.vanson.httpd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vanson.httpd</string>
    <key>ProgramArguments</key>
    <array>
        <string>/var/jb/usr/local/bin/vansonmodd</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DYLD_INSERT_LIBRARIES</key>
        <string>/var/jb/usr/lib/ellekit/libsubstitute.dylib</string>
        <key>HOME</key>
        <string>/var/jb/var/mobile</string>
        <key>TMPDIR</key>
        <string>/var/jb/var/tmp</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <!-- KeepAlive = true (非 dict) 最可靠：任何原因退出都会重启 -->
    <key>KeepAlive</key>
    <true/>
    <!-- 拒绝内存压力退出请求 -->
    <key>EnablePressuredExit</key>
    <false/>
    <!-- 进程调度优先级 -->
    <key>Nice</key>
    <integer>-10</integer>
    <key>LowPriorityIO</key>
    <false/>
    <!-- launchd 重启间隔(秒)：防止连续崩溃被限流 -->
    <key>ThrottleInterval</key>
    <integer>3</integer>
    <!-- 退出超时 -->
    <key>ExitTimeOut</key>
    <integer>10</integer>
    <!-- 启用事务追踪，launchd 知道进程是活跃的 -->
    <key>EnableTransactions</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/jb/var/log/vansonmodd.log</string>
    <key>StandardErrorPath</key>
    <string>/var/jb/var/log/vansonmodd_err.log</string>
</dict>
</plist>
