# daemon/com.vanson.httpd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vanson.httpd</string>
    <key>ProgramArguments</key>
    <array>
        <string>/var/jb/usr/local/bin/vansonmodd</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DYLD_INSERT_LIBRARIES</key>
        <string>/var/jb/usr/lib/ellekit/libsubstitute.dylib</string>
        <key>HOME</key>
        <string>/var/jb/var/mobile</string>
        <key>TMPDIR</key>
        <string>/var/jb/var/tmp</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>NetworkState</key>
        <true/>
    </dict>
    <!-- 关键: 不加 ProcessType → 默认 Interactive 优先级，Jetsam 不轻易杀 -->
    <!-- 不加 POSIXSpawnType → 不用 Adaptive，防止 launchd 自行管理生命周期 -->
    <key>EnablePressuredExit</key>
    <false/>
    <key>Nice</key>
    <integer>-5</integer>
    <key>LowPriorityIO</key>
    <false/>
    <key>ExitTimeOut</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/var/jb/var/log/vansonmodd.log</string>
    <key>StandardErrorPath</key>
    <string>/var/jb/var/log/vansonmodd_err.log</string>
</dict>
</plist>
